# etc/periodic/daily/apk-autoupgrade
#!/bin/sh
# I decided this was a good idea, it's probably not.

apk upgrade --update | ts '[%Y-%m-%d %H:%M:%S]' >> /var/log/apk-autoupgrade.log
